# deploy/systemd/small-talk-app.service
[Unit]
Description=Small Talk - App Server
Wants=network-online.target
After=network-online.target

[Service]
# Adjust ownership/path if you deploy differently
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/small-talk/apps/app/server
ExecStart=/home/ubuntu/small-talk/apps/app/server/app
EnvironmentFile=/home/ubuntu/small-talk/apps/app/server/app.env
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20
NoNewPrivileges=true

# Journald log identifier
SyslogIdentifier=small-talk-app

[Install]
WantedBy=multi-user.target


